# docker/testing-host/etc/s6-overlay/s6-rc.d/ssh/run
#!/command/execlineb -P

# Start SSH daemon
/usr/sbin/sshd -D -o ListenAddress=0.0.0.0 -o Port=22
