# docker/testing-host/etc/s6-overlay/s6-rc.d/coolify-network/up
#!/command/execlineb -P

# Create the coolify network if it doesn't exist
foreground { bash /etc/s6-overlay/scripts/create-coolify-network.sh }
