# docker/development/etc/s6-overlay/s6-rc.d/init-setup/up
#!/command/execlineb -P
foreground { composer -d /var/www/html/ install }
foreground { php /var/www/html/artisan migrate --step }
foreground { php /var/www/html/artisan dev --init }

