# docker/prod/etc/s6-overlay/s6-rc.d/init-script/up
#!/command/execlineb -P
s6-setuidgid www-data
php /var/www/html/artisan app:init
