# docker/prod/etc/s6-overlay/s6-rc.d/init-script/up
#!/command/execlineb -P
s6-setuidgid webuser
foreground {
    php /var/www/html/artisan optimize:clear
    php /var/www/html/artisan optimize
    php /var/www/html/artisan app:init
}
