# docker/prod-ssu/etc/s6-overlay/s6-rc.d/init-script/up
#!/command/execlineb -P
php /var/www/html/artisan app:init
