# docker/prod/etc/s6-overlay/s6-rc.d/init-seeder/up
#!/command/execlineb -P
php /var/www/html/artisan db:seed --class ProductionSeeder --force
