# docker/prod-ssu/etc/s6-overlay/s6-rc.d/ssh-agent/up
#!/usr/bin/execlineb -P
foreground {
  s6-sleep 5
  su - webuser -c "ssh-agent -a /tmp/coolify-ssh-agent.sock"
}
