# docker/testing-host/etc/s6-overlay/s6-rc.d/docker/run
#!/command/execlineb -P

# Start Docker daemon
/usr/bin/dockerd
