# docker/prod/etc/s6-overlay/s6-rc.d/scheduler-worker/run
#!/command/execlineb -P
foreground {
  s6-sleep 5
  su - webuser -c "php /var/www/html/artisan start:scheduler"
}
