# docker/dev-ssu/etc/s6-overlay/s6-rc.d/queue-worker/run
#!/command/execlineb -P
su - webuser -c "php /var/www/html/artisan queue:listen --timeout=600 --rest=1 --memory=512"
