# docker/prod-ssu/etc/s6-overlay/s6-rc.d/db-migration/up
#!/command/execlineb -P
php /var/www/html/artisan migrate --force --isolated
