# docker/prod-ssu/etc/s6-overlay/s6-rc.d/db-migration/up
#!/command/execlineb -P
VERSION=$(php bootstrap/getVersion.php)
echo "$VERSION"
# if [ echo "$VERSION" | grep -q "nightly" ]; then
#     php /var/www/html/artisan migrate:fresh --force
#     php /var/www/html/artisan db:seed --class ProductionSeeder --force
# else
php /var/www/html/artisan migrate --force --isolated
# fi
